# deploy/pedi-ai-app.service
[Unit]
Description=Pedi-AI App (Next.js)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/pedi-ai-app
ExecStart=/usr/bin/node .next/standalone/server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=/root/pedi-ai-app/.env

[Install]
WantedBy=multi-user.target